# matrix/continuwuity.service
[Unit]
Description=Continuwuity Matrix homeserver
Wants=network-online.target
After=network-online.target

[Service]
Type=notify
User=continuwuity
Group=continuwuity
ExecStart=/usr/local/bin/continuwuity --config /etc/continuwuity/continuwuity.toml
Restart=on-failure
RestartSec=5
# The first start after a Conduit migration can take several minutes before
# Continuwuity reports readiness to systemd.
TimeoutStartSec=15min

[Install]
WantedBy=multi-user.target
