# ssdnodes/includes/cert-renew.timer
[Unit]
Description=Try renew Lets Encrypt certificate twice a day

[Timer]
OnUnitActiveSec=12h
OnBootSec=1h

[Install]
WantedBy=timers.target
