# ssdnodes/includes/cert-renew.service
[Unit]
Description=Renew Lets Encrypt certificate

[Service]
Type=oneshot
ExecStart=/usr/bin/certbot renew
ExecStartPost=systemctl restart httpd
